# ruby/templates/etc/profile.d/rbenv.sh.erb
## THIS FILE IS MANAGED BY PUPPET
## modules<%= template_source.gsub(Regexp.new("^#{Puppet::Node::Environment.current[:modulepath].gsub(':','|')}"),"") %>

if [ `which rbenv` ]; then
    eval "$(rbenv init -)"
fi
