# supervisor/templates/etc/supervisor/conf.d/program.conf.erb
[program:<%= @program %>]
command = <%= @command %>
redirect_stderr = <%= @redirect_stderr %>
stdout_logfile = <%= @stdout_logfile %>

<%- if @user -%>
user = <%= @user %>
<%- end -%>

<%- if @directory -%>
directory = <%= @directory %>
<%- end -%>
