# supervisor/templates/etc/supervisor/conf.d/program.conf.erb
[program:<%= @program %>]
command = <%= @command %>
<%- if @user -%>
    user = <%= @user %>
<%- end -%>
